# systemd/geosite-update.service
[Unit]
Description=Update Xray geosite.dat
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-geosite.sh
